# etc/inittab
# --- T2-COPYRIGHT-NOTE-BEGIN ---
# This copyright note is auto-generated by ./scripts/Create-CopyPatch.
# 
# T2 SDE: package/.../sysfiles/etc_inittab.txt
# Copyright (C) 2004 - 2008 The T2 SDE Project
# Copyright (C) 1998 - 2003 ROCK Linux Project
# 
# More information can be found in the files COPYING and README.
# 
# This program is free software; you can redistribute it and/or modify
# it under the terms of the GNU General Public License as published by
# the Free Software Foundation; version 2 of the License. A copy of the
# GNU General Public License can be found in the file COPYING.
# --- T2-COPYRIGHT-NOTE-END ---
#
# See "man 5 inittab" for details.
#

# The default runlevel.
id:2:initdefault:

# /etc/rc.d executes the S and K scripts upon change
# of runlevel.
#
# Runlevels 0 and 6 are halt and reboot.
# Runlevels 1-5 are multi-user.
# Runlevel S is single user.
l0:0:wait:/etc/rc.d/rc
l1:1:wait:/etc/rc.d/rc
l2:2:wait:/etc/rc.d/rc
l3:3:wait:/etc/rc.d/rc
l4:4:wait:/etc/rc.d/rc
l5:5:wait:/etc/rc.d/rc
l6:6:wait:/etc/rc.d/rc
ls:S:wait:/etc/rc.d/rc
~~:S:wait:/sbin/sulogin

# What to do when CTRL-ALT-DEL is pressed.
ca:12345:ctrlaltdel:/sbin/shutdown -r now

# /sbin/agetty invocations for the runlevels.
#
# The "id" field MUST be the same as the last
# characters of the device (after "tty").
#
# Format:
#  <id>:<runlevels>:<action>:<process>

#
# Maintenace Console on /dev/console
#
#1:12345:respawn:/sbin/agetty -L -i -I '\012\015\012Maintenance Console:\012' 9600 console linux

#
# Virtual Consoles 1-6
#
# to enable the SAK blocking feature add the options "-t 10 -s /etc/sakmsg"
#
# to clear the screen on ANSI terminals add the option "-I '\033[H\033[J'"
#
# there are 3 Issue files:
#   /etc/issue ........ Simple login screen for consoles or terminals
#   /etc/issue.ansi ... Login screen for ANSI consoles or terminals
#   /etc/issue.net .... Login message for network connections
#
1:12345:respawn:/sbin/agetty -L -i -I '\012\015' 38400 vc/1 linux
2:12345:respawn:/sbin/agetty -f /etc/issue.ansi 38400 vc/2 linux
3:12345:respawn:/sbin/agetty -f /etc/issue.ansi 38400 vc/3 linux
4:12345:respawn:/sbin/agetty -f /etc/issue.ansi 38400 vc/4 linux
5:12345:respawn:/sbin/agetty -f /etc/issue.ansi 38400 vc/5 linux
6:12345:respawn:/sbin/agetty -f /etc/issue.ansi 38400 vc/6 linux

#T0:2345:respawn:/sbin/agetty   115200 ttyS0 vt100
#
# A X display manager (multiplexer)
#
X:5:respawn:/usr/X11/bin/startxdm 

