# etc/profile.d/less

export LESS="--ignore-case"
export LESSOPEN="|lesspipe.sh %s"

if [ -f "" ] ; then
LESSOPEN="| %s"
else
LESSOPEN="|lesspipe %s"
fi
export LESSOPEN
